# systemd/amsf-backup.service
[Unit]
Description=AMSF Google Drive database backup
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=guna999
Group=guna999
WorkingDirectory=/opt/apps/amsf-app
ExecStart=/opt/apps/amsf-app/.venv/bin/python /opt/apps/amsf-app/backup_db_to_google_drive.py